# deploy/sync.service
[Unit]
Description=Sync git pull (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/Sync
ExecStart=/var/www/Sync/deploy/sync.sh
Nice=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
